# tado-x-exporter/rootfs/etc/services.d/tado_x_exporter/run
#!/usr/bin/with-contenv bashio

declare tado_username
declare tado_password
declare tado_token_file_path
declare tado_exporter_port
declare tado_exporter_refresh_rate

TADO_USERNAME=$(bashio::config 'tado_username')
TADO_PASSWORD=$(bashio::config 'tado_password')
TADO_TOKEN_FILE_PATH=$(bashio::config 'tado_token_file_path')
TADO_EXPORTER_PORT=$(bashio::config 'tado_exporter_port')
TADO_EXPORTER_REFRESH_RATE=$(bashio::config 'tado_exporter_refresh_rate')

export TADO_USERNAME
export TADO_PASSWORD
export TADO_TOKEN_FILE_PATH
export TADO_EXPORTER_PORT
export TADO_EXPORTER_REFRESH_RATE

#bashio::log.info "${TADO_EXPORTER_PORT:="asdasdasd"}"
exec /usr/bin/tado_x_exporter